# deploy/ansible/roles/praho-native/templates/praho-qcluster.service.j2
# PRAHO Q-Cluster Service (Django-Q2 background tasks)
# Generated by Ansible - DO NOT EDIT

[Unit]
Description=PRAHO Q-Cluster (Django-Q2 background worker)
After=postgresql.service praho-platform.service
Wants=praho-platform.service

[Service]
Type=simple
User={{ project_user }}
Group={{ project_group }}
WorkingDirectory={{ project_root }}/src/services/platform
Environment=PYTHONPATH={{ project_root }}/src/services/platform
EnvironmentFile={{ project_root }}/.env
ExecStart={{ project_root }}/.venv/bin/python manage.py qcluster
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax={{ qcluster_memory_max }}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ project_root }}/logs {{ project_root }}/media

[Install]
WantedBy=multi-user.target
