# deploy/ansible/roles/praho-native/templates/praho-portal.service.j2
# PRAHO Portal Service (Native)
# Generated by Ansible - DO NOT EDIT

[Unit]
Description=PRAHO Portal (Django customer-facing service)
After=praho-platform.service network.target
Wants=praho-platform.service

[Service]
Type=notify
User={{ project_user }}
Group={{ project_group }}
WorkingDirectory={{ project_root }}/src/services/portal
EnvironmentFile={{ project_root }}/.env
ExecStart={{ project_root }}/.venv/bin/gunicorn \
    --bind 127.0.0.1:{{ portal_port }} \
    --workers {{ gunicorn_workers_portal }} \
    --timeout {{ gunicorn_timeout_portal }} \
    --access-logfile - \
    --error-logfile - \
    config.wsgi:application
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax={{ portal_memory_max }}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ project_root }}/logs

[Install]
WantedBy=multi-user.target
