# deploy/ansible/roles/praho-native/templates/praho-platform.service.j2
# PRAHO Platform Service (Native)
# Generated by Ansible - DO NOT EDIT

[Unit]
Description=PRAHO Platform (Django staff/admin service)
After=postgresql.service network.target
Requires=postgresql.service

[Service]
Type=notify
User={{ project_user }}
Group={{ project_group }}
WorkingDirectory={{ project_root }}/src/services/platform
Environment=PYTHONPATH={{ project_root }}/src/services/platform
EnvironmentFile={{ project_root }}/.env
ExecStart={{ project_root }}/.venv/bin/gunicorn \
    --bind 127.0.0.1:{{ platform_port }} \
    --workers {{ gunicorn_workers_platform }} \
    --timeout {{ gunicorn_timeout_platform }} \
    --access-logfile - \
    --error-logfile - \
    config.wsgi:application
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax={{ platform_memory_max }}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ project_root }}/static {{ project_root }}/media {{ project_root }}/logs

[Install]
WantedBy=multi-user.target
